# overlays/resize2fs/etc/rc.local
#!/bin/bash

/sbin/resize2fs `blkid -L rootfs`
